# greeting-mcp-node/systemd/mcp-greeting.service
[Unit]
Description=ScaleKit Greeting MCP Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/mcp-auth-demos/greeting-mcp-node
EnvironmentFile=/root/mcp-auth-demos/greeting-mcp-node/.env
ExecStart=/usr/bin/node dist/main.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
