# greeting-mcp-node/systemd/mcp-greeting-watch.service
[Unit]
Description=Restart MCP greeting server after rebuild

[Service]
Type=oneshot
ExecStart=systemctl restart mcp-greeting
